# deploy/pocketbase-backup.service
[Unit]
Description=PocketBase daily SQLite backup

[Service]
Type=oneshot
User=pocketbase
Group=pocketbase
ExecStart=/opt/pocketbase/backup-pocketbase.sh
